# units/bone-capemgr.target
[Unit]
Description=Device Tree Overlays
Documentation=https://github.com/hashrabbit/systemd-beaglebone#readme
